# glimmer.service.example
[Unit]
Description=Glimmer LED Controller
After=network.target

[Service]
Type=simple
# Run as root to access /dev/mem for GPIO control (rpi_ws281x requires root)
WorkingDirectory=/home/liquenops/Glimmer
# Use Python from virtual environment
ExecStart=/home/liquenops/Glimmer/venv/bin/python3 /home/liquenops/Glimmer/main.py --api-port 1129 --audio-port 31337
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
