# glimmer.service.example
[Unit]
Description=Glimmer LED Controller
After=network.target

[Service]
Type=simple
User=liquenops
WorkingDirectory=/home/liquenops/Glimmer
ExecStart=/usr/bin/python3 /home/liquenops/Glimmer/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# For simulator mode, use:
# ExecStart=/usr/bin/python3 /home/liquenops/Glimmer/main.py --simulator

# For custom options:
# ExecStart=/usr/bin/python3 /home/liquenops/Glimmer/main.py --api-port 1129 --audio-port 31337

[Install]
WantedBy=multi-user.target
